# roles/worker/templates/redis.service.j2
[Unit]
Description=Redis Server - Job Queue Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/lib/redis
ExecStart=/usr/bin/redis-server /etc/redis/redis.conf \
  --protected-mode no \
  --bind 127.0.0.1
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target